# overlay/etc/systemd/system/piuio2key.service
[Unit]
Description=PIUIO2Key-Linux USB IO bridge
Documentation=https://github.com/carlos-garcia/PIUIO2Key-Linux
After=local-fs.target
# NOT enabled by default — XSanity reads PIUIO natively.
# piu-launch.sh starts/stops this for pumptools keyboard-IO mode only.

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/bootinsanity/piuio2key/piu_bridge.py --no-tray
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
# WantedBy intentionally omitted — do not enable at boot.
